# docker/night-trigger/crontab
0 1 * * * root /apps/scripts/jenkins.sh
